# init.d/monet-dhcpc
#!/sbin/openrc-run
# -*- coding: utf-8; mode: sh; tab-width: 3; indent-tabs-mode: nil -*-
#
# Copyright 2018-2024 Raffaello D. Di Napoli
# Distributed under the terms of the GNU General Public License v2

prodname=${RC_SVCNAME%%-*}
. /lib/${prodname}/shared.shlib

svcitf="${RC_SVCNAME#*.}"
[ "${svcitf}" != ${RC_SVCNAME} ] || svcitf=

eval dhcpc_daemon=\${${svcitf}_dhcpc_daemon}
if [ -n "${dhcpc_daemon}" ]; then
   # TODO: discover what software it is.
   :
elif type dhcpcd >/dev/null 2>&1; then
   dhcpc_daemon=dhcpcd
elif busybox 2>&1 | grep -q '\sudhcpc,'; then
   dhcpc_daemon=udhcpc
fi

depend() {
   use logger
   # after because we can start without a hostname; none will be provided to
   # the DHCP server.
   after hostname
   # Not need, because we’ll mark it as inactive when we don’t have a lease,
   # and that would stop this service and never get a lease again.
   # Not use, because we do want OpenRC to start it.
   want ${svcitf:+${prodname}.${svcitf}}
}

if [ "${dhcpc_daemon##*/}" = dhcpcd ]; then
   name="DHCP client for ${svcitf} using dhcpcd"
   command="${dhcpc_daemon}"
   pidfile=/run/dhcpcd/${svcitf}.pid
   conffile=/etc/${prodname}/dhcpcd.${svcitf}.conf
   command_args="
      --config ${conffile} --script /lib/${prodname}/dhcpcd-event
      --background --noconfigure --nodev --nohook resolv.conf --nolink
   "
   if [ -n "$(conf_get_master ${svcitf})" ]; then
      # In the script, we assign the leased address to the bridge instead of
      # ${svcitf}; this goes wrong when dhcpcd validates the address with ARP,
      # because it finds the IP associated to the bridge’s MAC address instead
      # of ${svcitf}’s, and warns about DAD, then abandons the lease and
      # obtains a new address.
      # A hacky fix for this is to just disable all ARP messaging, which
      # results in DAD not happening.
      command_args="${command_args} --noarp"
   fi
   command_args="${command_args}
      --option domain_name,domain_name_servers,domain_search
      --option classless_static_routes,interface_mtu,rapid_commit
      ${svcitf}
   "
elif [ "${dhcpc_daemon}" = udhcpc ]; then
   name="DHCP client for ${svcitf} using udhcpc"
   command=busybox
   procname=udhcpc
   pidfile=/run/${prodname}/udhcpc.${svcitf}.pid
   command_args="
      -RS -A 5 -b -i ${svcitf} -p ${pidfile}
      -s /lib/${prodname}/udhcpc-event
   "
   hn=$(hostname)
   [ "x${hn}" != 'x(unknown)' ] || hn=
   if [ -n "${hn}" ]; then
      command_args="${command_args} -x hostname:${hn}"
   fi
   unset hn
elif [ -n "${svcitf}" ]; then
   name="DHCP client for ${svcitf}"
else
   name='DHCP client (make a symlink with suffix .interface)'
fi
description="${name}"

# Service implementation

start_pre() {
   if [ -z "${svcitf}" ]; then
      eerror "${RC_SVCNAME} cannot be started directly; you must create a"
      eerror "symlink to it named /etc/init.d/${RC_SVCNAME}.itf, where itf"
      eerror 'is the name of the network interface that you want to get an'
      eerror 'address via DHCP.'
      eerror ''
      eerror 'If you then enable rc_hotplug for it, this service will be'
      eerror 'started any time the link on the interface is up.'
      return 1
   fi

   if [ "${dhcpc_daemon##*/}" = dhcpcd ]; then
      mkdir -p /run/dhcpcd
   elif [ "${dhcpc_daemon}" = udhcpc ]; then
      mkdir -p /run/${prodname}
   fi
}
